# .tmp_ex/rootfs/etc/inittab
# BusyBox inittab
::sysinit:/etc/rc.init
::askfirst:-/bin/sh
::ctrlaltdel:/bin/reboot
::shutdown:/bin/umount -a -r
